# rc/.config/rc/init.d/rclone-nextcloud
#!/sbin/openrc-run


RCLONE_REMOTE=${RCLONE_REMOTE:-rackspace-nextcloud}
RCLONE_MOUNT_DIR=${RCLONE_MOUNT_DIR:-~/nextcloud}

supervisor=supervise-daemon
description="Mount with rclone (rackspace-nextcloud by default)"
command="/usr/bin/rclone"

command_args="mount \
  --vfs-cache-mode full \
  ${RCLONE_REMOTE}: \
  ${RCLONE_MOUNT_DIR}"

depend() {
  after net
}
